# roles/kopia/templates/kopia-backup.timer.j2
[Unit]
Description=Kopia Backup Timer

[Timer]
OnCalendar={{ kopia_backup_schedule }}
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
